# deployments/systemd/sunbreak.service
[Unit]
Description=Sunbreak keyword monitoring service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sunbreak
Group=sunbreak
WorkingDirectory=/opt/sunbreak
ExecStart=/usr/local/bin/sunbreak -config /etc/sunbreak/config.yaml
Restart=always
RestartSec=5
Environment=REDDIT_CLIENT_ID=
Environment=REDDIT_CLIENT_SECRET=
Environment=REDDIT_USER_AGENT=sunbreak-monitor/0.1
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/sunbreak

[Install]
WantedBy=multi-user.target
